# deploy/aigov-ingest.timer
[Unit]
Description=Run aigov ingest every hour

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=timers.target
